# dockerize/portalwatch_services/head/crontab
#30 * * * * . $INSTALL_PATH/$SCRIPT
0 15 * * 1 . $INSTALL_PATH/$SCRIPT
