# dockerize/portalwatch_services/datafetch/crontab
35 * * * * . $INSTALL_PATH/$SCRIPT
#0 2 * * 1 . $INSTALL_PATH/$SCRIPT
