# dockerize/portalwatch_services/datafetch/crontab
0 7 * * 1 . /adequate-portalharvester/datafetch.sh