# dockerize/portalwatch_services/datafetch/crontab
#0 7 * * 1 . /adequate-portalharvester/datafetch.sh
* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1