# dockerize/portalwatch_services/metafetch/crontab
0 4 * * 1 . /adequate-portalharvester/metafetch.sh

